# roles/adfc_permissions/files/adfc-fix-perms.service
[Unit]
Description=ADFC Fix Permissions

[Service]
User=root
ExecStart=/usr/local/sbin/fix_adfc_perms.sh